# telegram_userbot/systemd/telegram-userbot.service
[Unit]
Description=Telegram User Bot - Auto Upload System
After=network.target
Wants=network.target

[Service]
Type=simple
User=telegram-userbot
Group=telegram-userbot
WorkingDirectory=/opt/telegram-userbot
Environment=PATH=/opt/telegram-userbot/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/telegram-userbot/venv/bin/python main.py
Restart=always
RestartSec=10
KillMode=process
TimeoutStopSec=30

# Environment variables
EnvironmentFile=-/opt/telegram-userbot/.env

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/telegram-userbot
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-userbot

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target