# telegram_userbot/systemd/telegram-userbot.service
[Unit]
Description=Telegram User Bot - Production Service
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=%h/telegram_userbot
ExecStart=%h/telegram_userbot/run_userbot.sh
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStartSec=60
TimeoutStopSec=30

# Environment isolation
Environment=PATH=%h/telegram_userbot/.venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-%h/telegram_userbot/.env

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=%h/telegram_userbot
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-userbot

[Install]
WantedBy=multi-user.target